# systemd/pi-health-server.service
[Unit]
Description=Raspberry Pi health monitor service
After=local-fs.target

[Service]
Type=exec
ExecStart=/usr/local/bin/pi-health-server
User=root
Restart=always
RestartSec=2
Environment=PI_HEALTH_GPIO_CHIP=/dev/gpiochip0
Environment=PI_HEALTH_AC_GPIO6_LINE=6
Environment=PI_HEALTH_GPIO6_ACTIVE_LOW=0
Environment=PI_HEALTH_ACFAIL_GPIO21_LINE=21
Environment=PI_HEALTH_GPIO21_ACTIVE_LOW=0
Environment=PI_HEALTH_LOG_PATH=/var/log/pi-health/ac-power-events.log
Environment=PI_HEALTH_HTTP_HOST=0.0.0.0
Environment=PI_HEALTH_HTTP_PORT=12002
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
